# get-temp.timer
[Unit]
Description=Get Temperature & Humidity

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=timers.target
