# lcd-temperature.timer
[Unit]
Description=LCD Temperature Timer

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=timers.target
